# package/base-files/files/etc/hotplug.d/ieee80211/06-wifi-migrate
#!/bin/sh

[ "${ACTION}" = "add" ] || return

/usr/bin/ucode /usr/share/ucode/wifi/migrate.uc
